# etc/cron.d/axlprint
MAILTO=""
*/5     *       *       *       *       root	/axlprint/axlprint.sh
